# .tools/build/novnc/conf.d/cypress.conf
[program:cypress]
command = bash -c "cypress open --project /e2e/tests"
startsecs = 65
restart = always