# .tools/build/novnc/conf.d/websockify.conf
[program:websockify]
command=websockify --web /usr/share/novnc 5800 localhost:5900
autorestart=true
startsecs = 10
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
